# tincan/board/common/overlay/usr/lib/systemd/system/tincan-relay-watchdog.service
[Unit]
Description=Tincan system relay watchdog service
Requires=tincan-relay-watchdog.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/relay-watchdog

