# tincan/board/common/overlay/usr/lib/systemd/system/tincan-mgetty-watchdog.service
[Unit]
Description=Tincan system mgetty watchdog service
Requires=tincan-mgetty-watchdog.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/mgetty-watchdog

