# tincan/board/common/overlay/usr/lib/systemd/system/tincan-flag-watchdog.service
[Unit]
Description=Tincan system flag watchdog service
Requires=tincan-flag-watchdog.timer

[Service]
Type=simple
ExecStart=/usr/local/bin/flag-watchdog

