# tincan/board/common/overlay/usr/lib/systemd/system/tincan-flag-watchdog.timer
[Unit]
Description=Tincan flag watchdog timer

[Timer]
Unit=tincan-flag-watchdog.service
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
