# tincan/board/common/overlay/usr/lib/systemd/system/tincan-update.timer
[Unit]
Description=Tincan update timer

[Timer]
Unit=tincan-update.service
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
