# tincan/board/common/overlay/etc/systemd/system/postgresql.service.d/local.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/sh -c "if [ ! -f /var/lib/pgsql/PG_VERSION ]; then /usr/bin/pg_ctl initdb -D /var/lib/pgsql -o --locale=C; fi"

