# tincan/board/common/overlay/etc/systemd/system/supervisord.service.d/local.conf
[Service]
ExecStart=
ExecStart=/usr/bin/supervisord -n -c /etc/supervisord.conf -u root
RuntimeDirectory=supervisor
RuntimeDirectoryMode=755
